# scripts/systemd/portfolio-celery-beat.service
[Unit]
Description=Portfolio AI Celery Beat Scheduler
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/home/kasadis/portfolio-ai/backend
Environment="PATH=/home/kasadis/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
Environment="NUMBA_CACHE_DIR=%h/.cache/numba"
Environment="HF_HOME=%h/.cache/huggingface"
EnvironmentFile=-/home/kasadis/portfolio-ai/backend/.env
RuntimeDirectory=portfolio-ai
ExecStart=/home/kasadis/portfolio-ai/backend/.venv/bin/celery -A app.celery_app beat --loglevel=info
Restart=always
RestartSec=10
StandardOutput=append:/var/log/portfolio-ai/celery-beat.log
StandardError=append:/var/log/portfolio-ai/celery-beat-error.log

[Install]
WantedBy=default.target
