# scripts/systemd/portfolio-hatchet-worker.service
[Unit]
Description=Portfolio AI Hatchet Worker
After=network.target hatchet-engine.service

[Service]
Type=simple
WorkingDirectory=%h/portfolio-ai/backend
Environment="PATH=%h/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
Environment="NUMBA_CACHE_DIR=%h/.cache/numba"
Environment="HF_HOME=%h/.cache/huggingface"
EnvironmentFile=-%h/portfolio-ai/backend/.env
EnvironmentFile=-%h/.env.local
ExecStart=%h/portfolio-ai/backend/.venv/bin/python -m app.worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-hatchet-worker

KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
