# scripts/systemd/portfolio-backend.service
[Unit]
Description=Portfolio AI Backend (FastAPI)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=%h/portfolio-ai/backend
# Added %h/.local/bin for claude/gemini CLIs
Environment="PATH=%h/.local/bin:%h/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
Environment="DB_POOL_SIZE=3"
Environment="DB_MAX_OVERFLOW=2"
Environment="NUMBA_CACHE_DIR=%h/.cache/numba"
Environment="HF_HOME=%h/.cache/huggingface"
EnvironmentFile=-%h/portfolio-ai/backend/.env
ExecStart=%h/portfolio-ai/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-backend

[Install]
WantedBy=default.target
