# scripts/systemd/portfolio-celery.service
[Unit]
Description=Portfolio AI Celery Worker
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=portfolio-ai
Group=portfolio-ai
WorkingDirectory=/home/kasadis/portfolio-ai/backend
Environment="PATH=/home/kasadis/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/var/cache/portfolio-ai"
Environment="NUMBA_CACHE_DIR=/var/cache/portfolio-ai/numba"
Environment="HF_HOME=/var/cache/portfolio-ai/huggingface"
EnvironmentFile=-/home/kasadis/portfolio-ai/backend/.env
RuntimeDirectory=portfolio-ai
ExecStart=/home/kasadis/portfolio-ai/backend/.venv/bin/celery -A app.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=10
StandardOutput=append:/var/log/portfolio-ai/celery-worker.log
StandardError=append:/var/log/portfolio-ai/celery-worker-error.log

[Install]
WantedBy=multi-user.target
