# scripts/systemd/portfolio-celery-beat.service
[Unit]
Description=Portfolio AI Celery Beat Scheduler
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=%h/portfolio-ai/backend
Environment="PATH=%h/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=%h"
Environment="NUMBA_CACHE_DIR=%h/.cache/numba"
Environment="HF_HOME=%h/.cache/huggingface"
EnvironmentFile=-%h/portfolio-ai/backend/.env
ExecStart=%h/portfolio-ai/backend/.venv/bin/celery -A app.celery_app beat --loglevel=info
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio-beat

# ZOMBIE PREVENTION: Kill all child processes when stopping
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
