# scripts/systemd/portfolio-frontend.service
[Unit]
Description=Portfolio AI Frontend (Next.js)
After=network.target portfolio-backend.service
Wants=portfolio-backend.service

[Service]
Type=simple
User=kasadis
WorkingDirectory=/home/kasadis/portfolio-ai/frontend
Environment="PATH=/home/kasadis/.nvm/versions/node/v20.11.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=development"
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
