# scripts/portfolio-ai-celery-worker.service
[Unit]
Description=Portfolio AI Celery Worker
After=network.target redis.service

[Service]
Type=forking
User=kasadis
Group=kasadis
WorkingDirectory=/home/kasadis/portfolio-ai/backend
Environment="PATH=/home/kasadis/portfolio-ai/backend/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/kasadis/portfolio-ai/backend/.venv/bin/celery -A app.celery_app worker --concurrency=2 --loglevel=info --logfile=/var/log/portfolio-ai-celery-worker.log --pidfile=/var/run/portfolio-ai-celery-worker.pid --detach
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
